# cos-csi-mounter/install/cos-csi-mounter.service
[Unit]
Description=COS CSI mounter service
DefaultDependencies=no
Requires=network-online.target
After=network-online.target
Before=shutdown.target reboot.target halt.target
StartLimitIntervalSec=0

[Service]
Type=notify
ExecStart=/usr/local/bin/cos-csi-mounter-server
KillMode=process

# Restart Configuration - Always restart the service
Restart=always
RestartSec=10

# Watchdog Configuration - Service must send keepalive every 30 seconds
WatchdogSec=30
NotifyAccess=main

# Timeouts
TimeoutStartSec=60
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cos-csi-mounter

[Install]
WantedBy=multi-user.target
